# rustdesk/rootfs/etc/s6-overlay/s6-rc.d/init-rustdesk/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Rustdesk
# Rustdesk init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring Rustdesk..."
