# rustdesk/rootfs/etc/s6-overlay/s6-rc.d/rustdesk-hbbs/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: NRustdesk8N
# Runs Rustdesk script
# ==============================================================================

bashio::log.info "Starting Rustdesk hbbs..."

exec hbbs